# templates/etc/profile.d/nvm.sh.j2
export NVM_DIR="{{ nvm.root }}"

. {{ nvm.root }}/nvm.sh
