# templates/etc/profile.d/nvm.sh.j2
if shopt -q login_shell ; then
  export NVM_DIR="{{ nvm.root }}"

  . {{ nvm.root }}/nvm.sh
fi
